# etc/macos-launchd/syncthing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
	Reference template — the install.sh script generates this automatically.

	For manual setup, replace USERNAME with your macOS username (e.g. "jb"),
	copy to ~/Library/LaunchAgents/syncthing.plist, and run:

	    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/syncthing.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.syncthing.syncthing</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/USERNAME/local/bin/syncthing</string>
		<string>--no-browser</string>
		<string>--no-restart</string>
		<string>--no-upgrade</string>
	</array>

	<key>KeepAlive</key>
	<true/>

	<key>LowPriorityIO</key>
	<true/>

	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>/Users/USERNAME/Library/Logs/Syncthing.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/USERNAME/Library/Logs/Syncthing-Errors.log</string>
</dict>
</plist>
